# systemd_services/piclock_display.service
[Unit]
Description=pi_clock Disaplay data

[Service]
Type=idle
ExecStart=/root/pi_clock/piclock_display.py

[Install]
WantedBy=multi-user.target
